# base/mariadb-backup/files/mariadb-backup@.service
[Unit]
Description=%i mariadb backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mariadb-backup %i
Nice=10
User=root
Group=root
