# fstrim/files/fstrim.service
[Unit]
Description=Discard unused blocks using TRIM/UNMAP

[Service]
Type=oneshot
ExecStart=/sbin/fstrim -av
